# example/echoes.service
[Unit]
Description=Echoes service

[Service]
Type=forking
ExecStart=/bin/bash -c "/home/com.mutantcat.echoes-1.0.20240405-linux_x86 -port 9967 -notice_mod dingbot -token xxx -secret0 xxx -server_name 自建-服务器1 &"
KillMode=process
Restart=on-failure
RestartSec=3s

[Install]
WantedBy=multi-user.target